# systemd/user/batteryNotification.service
[Unit]
Description=Battery notifyer

[Service]
Type=oneshot
ExecStart=/home/jeanclaude/bin/battery.sh -c
