# systemd/user/mbsync.service
[Unit]
Description=Mailbox syncronisation service

[Service]
Type=oneshot
ExecStart=/home/jeanclaude/.config/systemd/user/mbsync.startup
