# systemd/user/mbsync.service
[Unit]
Description=Mailbox syncronisation service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Va
