# systemd/user/vdirsyncer.service
[Unit]
Description=DAV syncronisation service

[Service]
Type=oneshot
ExecStart=/home/jeanclaude/.config/systemd/user/vdirsyncer.startup
