# systemd/user/vdirsyncer.service
[Unit]
Description=DAV syncronisation service

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
