# CalCard/dav_syncer.service
[Unit]
Description=CalDAV and CardDAV syncer

[Service]
Type=oneshot
ExecStart=/home/jeanclaude/.config/systemd/user/dav_syncer.startup
